# cast/systemd/screencast.socket
[Unit]
Before=screencast.service

[Socket]
ListenStream=554
ListenStream=5554

[Install]
WantedBy=sockets.target
