# cast/systemd/screencast.socket
[Unit]
Before=screencast.service

[Socket]
ListenStream=554
Accept=no

[Install]
WantedBy=sockets.target
